# .config/systemd/user/udiskie.service
[Unit]
Description=udiskie automounter

[Service]
ExecStart=/usr/bin/udiskie -an

[Install]
WantedBy=default.target
